# ft_services/srcs/influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:influx]
command = sh -c "influxd && pkill supervisord"

[program:telegraf]
command = sh -c "/usr/telegraf/usr/bin/telegraf && pkill supervisord"